# .config/systemd/user/dwl.service
[Unit]
Description=dwl after login daemon
After=getty.target

[Service]
ExecStart=%h/.local/bin/dwl

[Install]
WantedBy=default.target
