# .config/systemd/user/gammastep.service
[Unit]
Description=gammastep daemon
After=dwl.service

[Service]
ExecStart=gammastep
Restart=on-failure

[Install]
WantedBy=default.target
